# systemd/weather-mon.socket
[Unit]
Description=Weather monitor socket

[Socket]
ListenStream=%t/weather-mon.sock

[Install]
WantedBy=sockets.target
